# debian/cron.d
*/30 * * * * root [ -x /usr/lib/architect/cron/contractorSync ] && /usr/lib/architect/cron/contractorSync
