# docker/etc/service/api_server/run
#!/bin/sh

exec /usr/local/architect/api_server/api_server.py
