# docker/etc/service/fake_monitor/run
#!/bin/sh

while true
do
 /usr/local/architect/demo_test/testData_monitor
 sleep 10
done
